# templates/default/hashi-ui.service.erb
## This file is dynamically generated by Chef. DO NOT EDIT manually.

[Unit]
Description=Hashi-ui
Documentation=https://github.com/jippi/hashi-ui
Wants=network-online.target
After=network-online.target

[Service]
User=<%= @user %>
Group=<%= @group %>

<% @options.each do |option_name,option_value| -%>
  <% if !!option_value == option_value %>
    <% @start_opt += " --#{option_name}" if option_value.is_a?(TrueClass) %>
  <% else %>
    <% @start_opt += " --#{option_name} #{option_value}" unless nil_or_empty?(option_value) %>
  <% end %>
<% end -%>
ExecStart=<%= @binary %> <%= @start_opt %>

StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target